# templates/openclaw-backup.service.j2
[Unit]
Description=Backup OpenClaw state and skills

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-openclaw.sh
